# root/etc/services.d/unbound/run
#!/usr/bin/execlineb -P

# not needed as unbound.conf is set to drop rights after launc (plus this fails if I change UID here coz it can't bind to port 53)
# s6-setuidgid unbound

/sbin/unbound -c /etc/unbound/unbound.conf
