# services.d/stubby/run
#!/usr/bin/execlineb -P

s6-setuidgid stubby

/usr/local/bin/stubby -C /usr/local/etc/stubby/stubby.yml