# etc/services.d/unbound/run
#!/usr/bin/execlineb -P

/usr/sbin/unbound-anchor -a /etc/unbound/root.key

# switch user
s6-setuidgid unbound
/usr/sbin/unbound -c /etc/unbound/unbound.conf
