# etc/services.d/unbound/run
#!/usr/bin/execlineb -P

s6-setuidgid unbound

/usr/sbin/unbound-anchor -a /etc/unbound/root.key
/usr/sbin/unbound -c /etc/unbound/unbound.conf